# image/sysroot-tree/usr/lib/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
[Unit]
Description=autologin
ConditionPathExists=/proc/cmdline
ConditionKernelCommandLine=|constellation.console
ConditionKernelCommandLine=|constellation.debug

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --keep-baud --autologin root 115200,57600,38400,9600 - $TERM

[Install]
WantedBy=multi-user.target
