# image/sysroot-tree/etc/systemd/system/sshd-keygen@.service.d/override.conf
[Unit]
ConditionFileNotEmpty=|!/var/run/state/ssh/ssh_host_%i_key
