# image/sysroot-tree/etc/systemd/system/sshd-keygen@.service.d/override.conf
[Unit]
ConditionFileNotEmpty=|!/var/run/state/ssh/ssh_host_%i_key

[Service]
ExecStart=/usr/bin/bash -c "/usr/bin/ssh-keygen -q -t %i -f /var/run/state/ssh/ssh_host_%i_key -C '' -N ''"
