# image/base/mkosi.skeleton/usr/lib/systemd/system/containerd.service.d/local.conf
[Service]
ExecStart=
ExecStart=/usr/bin/containerd --config /usr/etc/containerd/config.toml

# Until https://github.com/containerd/containerd/pull/8924 lands in our
# containerd version, we need to decrease the default ulimit that
# is set in the upstream containerd service of infinity which newer
# systemd versions resolve to 2**30 which is way too large and
# results in various inspecific errors such as excessive resource (e.g., memory)
# usage.
LimitNOFILE=524288
