# image/sysroot-tree/usr/lib/systemd/system/sshd-keygen.target
[Unit]
Wants=sshd-keygen@ed25519.service
PartOf=sshd.service
