# image/sysroot-tree/usr/lib/systemd/system/sshd-keygen.target
# Disable automatic key generation
[Unit]
PartOf=sshd.service
