# mkosi/mkosi.skeleton/usr/lib/systemd/system/containerd.service.d/local.conf
[Service]
ExecStart=
ExecStart=/usr/bin/containerd --config /usr/etc/containerd/config.toml
