# image/sysroot-tree/usr/lib/systemd/system/getty@tty1.service.d/autologin.conf
[Unit]
Description=autologin
ConditionPathExists=/proc/cmdline
ConditionKernelCommandLine=|constel.console
ConditionKernelCommandLine=|constel.debug

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root %I $TERM
[Install]
WantedBy=multi-user.target
