# image/mkosi/mkosi.skeleton/etc/profile.d/constellation.sh
#!/usr/bin/env bash
export PATH=/run/state/bin:${PATH}
export KUBECONFIG=/etc/kubernetes/admin.conf
alias k=kubectl
